# emulationstation/service/DEBIAN/postrm
#!/bin/sh

set -e

update-rc.d web2rgbmatrix remove

exit 0
