# emulationstation/service/DEBIAN/postinst
#!/bin/sh

set -e 

update-rc.d web2rgbmatrix defaults
service web2rgbmatrix start
chown -R pi:pi /etc/web2rgbmatrix
chown -R pi:pi /home/pi/.emulationstation/scripts
exit 0
